# src/main/config/cron.d/check-cdp-listend.cron
#  cron to check if ${project.artifactId} daemon is alive and restart it if needed
#
0 * * * * root (date --iso-8601=seconds; PATH=/sbin:/bin:/usr/sbin:/usr/bin; /usr/quattor/scripts/check-${project.artifactId}) >> /var/log/check-${project.artifactId}.cron.log 2>&1

